# config/supervisord/rsyslog.ini
[program:rsyslog]
command=/sbin/rsyslogd -n
numprocs=1
autostart=true
autorestart=true
